# projects/ROCKNIX/packages/tools/rocknix-splash/system.d/rocknix-splash.service
[Unit]
Description=ROCKNIX Splash Screen
Requires=tmp.mount
After=tmp.mount

[Service]
Type=oneshot
ExecStart=/usr/bin/rocknix-splash
ExecStartPost=/bin/sh -c "echo -en \"\033[1000H\033[2K$(cat /etc/issue)\n\n\" >/dev/console && echo -ne \"\033[1000H\033[2K==> Loading, please wait...\" >/dev/console"

[Install]
WantedBy=sysinit.target
